# resources/supervisord/redis.sample.conf
[program:redis]
command=redis-server --port 3001
stdout_logfile=/var/log/supervisor/redis.err.log
redirect_stderr=true
autostart=true
autorestart=true
