# resources/supervisor/queues.conf
[program:QUEUE_CHART]
process_name=%(program_name)s_%(process_num)02d
directory=APP_CURRENT_PATH
command=php artisan queue:work beanstalkd --force --queue=QUEUE_CHART --tries=3
redirect_stderr=0
stdout_logfile=/var/log/supervisor/SUPERVISOR_GROUP/%(program_name)s_stdout.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/SUPERVISOR_GROUP/%(program_name)s_stderr.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=3
autostart=1
autorestart=1
user=APP_SERVER_USER
numprocs=1

[program:QUEUE_CLASSIFICATION]
process_name=%(program_name)s_%(process_num)02d
directory=APP_CURRENT_PATH
command=php artisan queue:work beanstalkd --force --queue=QUEUE_CLASSIFICATION --tries=3
redirect_stderr=0
stdout_logfile=/var/log/supervisor/SUPERVISOR_GROUP/%(program_name)s_stdout.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/SUPERVISOR_GROUP/%(program_name)s_stderr.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=3
autostart=1
autorestart=1
user=APP_SERVER_USER
numprocs=1

[program:QUEUE_DEFAULT]
process_name=%(program_name)s_%(process_num)02d
directory=APP_CURRENT_PATH
command=php artisan queue:work beanstalkd --force --queue=QUEUE_DEFAULT --tries=3
redirect_stderr=0
stdout_logfile=/var/log/supervisor/SUPERVISOR_GROUP/%(program_name)s_stdout.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/SUPERVISOR_GROUP/%(program_name)s_stderr.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=3
autostart=1
autorestart=1
user=APP_SERVER_USER
numprocs=1

[program:QUEUE_EVENT]
process_name=%(program_name)s_%(process_num)02d
directory=APP_CURRENT_PATH
command=php artisan queue:work beanstalkd --force --queue=QUEUE_EVENT --tries=3
redirect_stderr=0
stdout_logfile=/var/log/supervisor/SUPERVISOR_GROUP/%(program_name)s_stdout.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/SUPERVISOR_GROUP/%(program_name)s_stderr.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=3
autostart=1
autorestart=1
user=APP_SERVER_USER
numprocs=1

[program:QUEUE_EXPORT]
process_name=%(program_name)s_%(process_num)02d
directory=APP_CURRENT_PATH
command=php artisan queue:work beanstalkd --force --queue=QUEUE_EXPORT --tries=3
redirect_stderr=0
stdout_logfile=/var/log/supervisor/SUPERVISOR_GROUP/%(program_name)s_stdout.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/SUPERVISOR_GROUP/%(program_name)s_stderr.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=3
autostart=1
autorestart=1
user=APP_SERVER_USER
numprocs=3
sleep=1

[program:QUEUE_GEOLOCATE]
process_name=%(program_name)s_%(process_num)02d
directory=APP_CURRENT_PATH
command=php artisan queue:work beanstalkd --force --queue=QUEUE_GEOLOCATE --tries=3
redirect_stderr=0
stdout_logfile=/var/log/supervisor/SUPERVISOR_GROUP/%(program_name)s_stdout.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/SUPERVISOR_GROUP/%(program_name)s_stderr.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=3
autostart=1
autorestart=1
user=APP_SERVER_USER
numprocs=1

[program:QUEUE_IMPORT]
process_name=%(program_name)s_%(process_num)02d
directory=APP_CURRENT_PATH
command=php artisan queue:work beanstalkd --force --queue=QUEUE_IMPORT --tries=3
redirect_stderr=0
stdout_logfile=/var/log/supervisor/SUPERVISOR_GROUP/%(program_name)s_stdout.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/SUPERVISOR_GROUP/%(program_name)s_stderr.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=3
autostart=1
autorestart=1
user=APP_SERVER_USER
numprocs=1

[program:QUEUE_OCR]
process_name=%(program_name)s_%(process_num)02d
directory=APP_CURRENT_PATH
command=php artisan queue:work beanstalkd --force --queue=QUEUE_OCR --tries=3
redirect_stderr=0
stdout_logfile=/var/log/supervisor/SUPERVISOR_GROUP/%(program_name)s_stdout.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/SUPERVISOR_GROUP/%(program_name)s_stderr.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=3
autostart=1
autorestart=1
user=APP_SERVER_USER
numprocs=1

[program:QUEUE_LAMBDA_OCR]
process_name=%(program_name)s_%(process_num)02d
directory=APP_CURRENT_PATH
command=php artisan queue:work beanstalkd --force --queue=QUEUE_LAMBDA_OCR --tries=3
redirect_stderr=0
stdout_logfile=/var/log/supervisor/SUPERVISOR_GROUP/%(program_name)s_stdout.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/SUPERVISOR_GROUP/%(program_name)s_stderr.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=3
autostart=1
autorestart=1
user=APP_SERVER_USER
numprocs=1

[program:QUEUE_BIOSPEX_EVENT]
process_name=%(program_name)s_%(process_num)02d
directory=APP_CURRENT_PATH
command=php artisan queue:work beanstalkd --force --queue=QUEUE_BIOSPEX_EVENT --tries=3
redirect_stderr=0
stdout_logfile=/var/log/supervisor/SUPERVISOR_GROUP/%(program_name)s_stdout.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/SUPERVISOR_GROUP/%(program_name)s_stderr.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=3
autostart=1
autorestart=1
user=APP_SERVER_USER
numprocs=1

[program:QUEUE_PUSHER_HANDLER]
process_name=%(program_name)s_%(process_num)02d
directory=APP_CURRENT_PATH
command=php artisan queue:work beanstalkd --force --queue=QUEUE_PUSHER_HANDLER --tries=3
redirect_stderr=0
stdout_logfile=/var/log/supervisor/SUPERVISOR_GROUP/%(program_name)s_stdout.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/SUPERVISOR_GROUP/%(program_name)s_stderr.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=3
autostart=1
autorestart=1
user=APP_SERVER_USER
numprocs=1

[program:QUEUE_WEDIGBIO_EVENT]
process_name=%(program_name)s_%(process_num)02d
directory=APP_CURRENT_PATH
command=php artisan queue:work beanstalkd --force --queue=QUEUE_WEDIGBIO_EVENT --tries=3
redirect_stderr=0
stdout_logfile=/var/log/supervisor/SUPERVISOR_GROUP/%(program_name)s_stdout.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/SUPERVISOR_GROUP/%(program_name)s_stderr.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=3
autostart=1
autorestart=1
user=APP_SERVER_USER
numprocs=1

[program:QUEUE_PUSHER_PROCESS]
process_name=%(program_name)s_%(process_num)02d
directory=APP_CURRENT_PATH
command=php artisan queue:work beanstalkd --force --queue=QUEUE_PUSHER_PROCESS --tries=3
redirect_stderr=0
stdout_logfile=/var/log/supervisor/SUPERVISOR_GROUP/%(program_name)s_stdout.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/SUPERVISOR_GROUP/%(program_name)s_stderr.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=3
autostart=1
autorestart=1
user=APP_SERVER_USER
numprocs=1

[program:QUEUE_RECONCILE]
process_name=%(program_name)s_%(process_num)02d
directory=APP_CURRENT_PATH
command=php artisan queue:work beanstalkd --force --queue=QUEUE_RECONCILE --tries=3
redirect_stderr=0
stdout_logfile=/var/log/supervisor/SUPERVISOR_GROUP/%(program_name)s_stdout.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/SUPERVISOR_GROUP/%(program_name)s_stderr.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=3
autostart=1
autorestart=1
user=APP_SERVER_USER
numprocs=1

[program:QUEUE_SNS_RECONCILIATION]
process_name=%(program_name)s_%(process_num)02d
directory=APP_CURRENT_PATH
command=php artisan queue:work beanstalkd --force --queue=QUEUE_SNS_RECONCILIATION --tries=3
redirect_stderr=0
stdout_logfile=/var/log/supervisor/SUPERVISOR_GROUP/%(program_name)s_stdout.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/SUPERVISOR_GROUP/%(program_name)s_stderr.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=3
autostart=1
autorestart=1
user=APP_SERVER_USER
numprocs=1

[program:QUEUE_SNS_TESSERACT_OCR]
process_name=%(program_name)s_%(process_num)02d
directory=APP_CURRENT_PATH
command=php artisan queue:work beanstalkd --force --queue=QUEUE_SNS_TESSERACT_OCR --tries=3
redirect_stderr=0
stdout_logfile=/var/log/supervisor/SUPERVISOR_GROUP/%(program_name)s_stdout.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/SUPERVISOR_GROUP/%(program_name)s_stderr.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=3
autostart=1
autorestart=1
user=APP_SERVER_USER
numprocs=1

[program:QUEUE_WORKFLOW]
process_name=%(program_name)s_%(process_num)02d
directory=APP_CURRENT_PATH
command=php artisan queue:work beanstalkd --force --queue=QUEUE_WORKFLOW --tries=3
redirect_stderr=0
stdout_logfile=/var/log/supervisor/SUPERVISOR_GROUP/%(program_name)s_stdout.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/SUPERVISOR_GROUP/%(program_name)s_stderr.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=3
autostart=1
autorestart=1
user=APP_SERVER_USER
numprocs=1

[program:QUEUE_SERNEC_FILE]
process_name=%(program_name)s_%(process_num)02d
directory=APP_CURRENT_PATH
command=php artisan queue:work beanstalkd --force --queue=QUEUE_SERNEC_FILE --tries=3
redirect_stderr=0
stdout_logfile=/var/log/supervisor/SUPERVISOR_GROUP/%(program_name)s_stdout.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/SUPERVISOR_GROUP/%(program_name)s_stderr.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=3
autostart=1
autorestart=1
user=APP_SERVER_USER
numprocs=1

[program:QUEUE_SERNEC_ROW]
process_name=%(program_name)s_%(process_num)02d
directory=APP_CURRENT_PATH
command=php artisan queue:work beanstalkd --force --queue=QUEUE_SERNEC_ROW --tries=3
redirect_stderr=0
stdout_logfile=/var/log/supervisor/SUPERVISOR_GROUP/%(program_name)s_stdout.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/SUPERVISOR_GROUP/%(program_name)s_stderr.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=3
autostart=1
autorestart=1
user=APP_SERVER_USER
numprocs=2
