# resources/supervisor/queues.conf
[program:QUEUE_CHART]
process_name=%(program_name)s_%(process_num)02d
directory=APP_CURRENT_PATH
command=php artisan queue:work beanstalkd --force --queue=QUEUE_CHART --tries=3
redirect_stderr=true
stdout_logfile=/var/log/supervisor/APP_DOMAIN/QUEUE_CHART.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
autostart=true
autorestart=true
user=APP_SERVER_USER
numprocs=QUEUE_NUM_PROCS

[program:QUEUE_CLASSIFICATION]
process_name=%(program_name)s_%(process_num)02d
directory=APP_CURRENT_PATH
command=php artisan queue:work beanstalkd --force --queue=QUEUE_CLASSIFICATION --tries=3
redirect_stderr=true
stdout_logfile=/var/log/supervisor/APP_DOMAIN/QUEUE_CLASSIFICATION.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
autostart=true
autorestart=true
user=APP_SERVER_USER
numprocs=QUEUE_NUM_PROCS

[program:QUEUE_DEFAULT]
process_name=%(program_name)s_%(process_num)02d
directory=APP_CURRENT_PATH
command=php artisan queue:work beanstalkd --force --queue=QUEUE_DEFAULT --tries=3
redirect_stderr=true
stdout_logfile=/var/log/supervisor/APP_DOMAIN/QUEUE_DEFAULT.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
autostart=true
autorestart=true
user=APP_SERVER_USER
numprocs=QUEUE_NUM_PROCS

[program:QUEUE_EVENT]
process_name=%(program_name)s_%(process_num)02d
directory=APP_CURRENT_PATH
command=php artisan queue:work beanstalkd --force --queue=QUEUE_EVENT --tries=3
redirect_stderr=true
stdout_logfile=/var/log/supervisor/APP_DOMAIN/QUEUE_EVENT.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
autostart=true
autorestart=true
user=APP_SERVER_USER
numprocs=QUEUE_NUM_PROCS

[program:QUEUE_EXPORT]
process_name=%(program_name)s_%(process_num)02d
directory=APP_CURRENT_PATH
command=php artisan queue:work beanstalkd --force --queue=QUEUE_EXPORT --tries=3
redirect_stderr=true
stdout_logfile=/var/log/supervisor/APP_DOMAIN/QUEUE_EXPORT.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
autostart=true
autorestart=true
user=APP_SERVER_USER
numprocs=QUEUE_NUM_PROCS

[program:QUEUE_GEOLOCATE]
process_name=%(program_name)s_%(process_num)02d
directory=APP_CURRENT_PATH
command=php artisan queue:work beanstalkd --force --queue=QUEUE_GEOLOCATE --tries=3
redirect_stderr=true
stdout_logfile=/var/log/supervisor/APP_DOMAIN/QUEUE_GEOLOCATE.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
autostart=true
autorestart=true
user=APP_SERVER_USER
numprocs=QUEUE_NUM_PROCS

[program:QUEUE_IMPORT]
process_name=%(program_name)s_%(process_num)02d
directory=APP_CURRENT_PATH
command=php artisan queue:work beanstalkd --force --queue=QUEUE_IMPORT --tries=3
redirect_stderr=true
stdout_logfile=/var/log/supervisor/APP_DOMAIN/QUEUE_IMPORT.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
autostart=true
autorestart=true
user=APP_SERVER_USER
numprocs=QUEUE_NUM_PROCS

[program:QUEUE_LAMBDA_OCR]
process_name=%(program_name)s_%(process_num)02d
directory=APP_CURRENT_PATH
command=php artisan queue:work beanstalkd --force --queue=QUEUE_LAMBDA_OCR --tries=3
redirect_stderr=true
stdout_logfile=/var/log/supervisor/APP_DOMAIN/QUEUE_LAMBDA_OCR.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
autostart=true
autorestart=true
user=APP_SERVER_USER
numprocs=QUEUE_NUM_PROCS

[program:QUEUE_BIOSPEX_EVENT]
process_name=%(program_name)s_%(process_num)02d
directory=APP_CURRENT_PATH
command=php artisan queue:work beanstalkd --force --queue=QUEUE_BIOSPEX_EVENT --tries=3
redirect_stderr=true
stdout_logfile=/var/log/supervisor/APP_DOMAIN/QUEUE_BIOSPEX_EVENT.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
autostart=true
autorestart=true
user=APP_SERVER_USER
numprocs=QUEUE_NUM_PROCS

[program:QUEUE_PUSHER_CLASSIFICATION]
process_name=%(program_name)s_%(process_num)02d
directory=APP_CURRENT_PATH
command=php artisan queue:work beanstalkd --force --queue=QUEUE_PUSHER_CLASSIFICATION --tries=3
redirect_stderr=true
stdout_logfile=/var/log/supervisor/APP_DOMAIN/QUEUE_PUSHER_CLASSIFICATION.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
autostart=true
autorestart=true
user=APP_SERVER_USER
numprocs=QUEUE_NUM_PROCS

[program:QUEUE_PUSHER_HANDLER]
process_name=%(program_name)s_%(process_num)02d
directory=APP_CURRENT_PATH
command=php artisan queue:work beanstalkd --force --queue=QUEUE_PUSHER_HANDLER --tries=3
redirect_stderr=true
stdout_logfile=/var/log/supervisor/APP_DOMAIN/QUEUE_PUSHER_HANDLER.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
autostart=true
autorestart=true
user=APP_SERVER_USER
numprocs=QUEUE_NUM_PROCS

[program:QUEUE_WEDIGBIO_EVENT]
process_name=%(program_name)s_%(process_num)02d
directory=APP_CURRENT_PATH
command=php artisan queue:work beanstalkd --force --queue=QUEUE_WEDIGBIO_EVENT --tries=3
redirect_stderr=true
stdout_logfile=/var/log/supervisor/APP_DOMAIN/QUEUE_WEDIGBIO_EVENT.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
autostart=true
autorestart=true
user=APP_SERVER_USER
numprocs=QUEUE_NUM_PROCS

[program:QUEUE_PUSHER_PROCESS]
process_name=%(program_name)s_%(process_num)02d
directory=APP_CURRENT_PATH
command=php artisan queue:work beanstalkd --force --queue=QUEUE_PUSHER_PROCESS --tries=3
redirect_stderr=true
stdout_logfile=/var/log/supervisor/APP_DOMAIN/QUEUE_PUSHER_PROCESS.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
autostart=true
autorestart=true
user=APP_SERVER_USER
numprocs=QUEUE_NUM_PROCS

[program:QUEUE_RECONCILE]
process_name=%(program_name)s_%(process_num)02d
directory=APP_CURRENT_PATH
command=php artisan queue:work beanstalkd --force --queue=QUEUE_RECONCILE --tries=3
redirect_stderr=true
stdout_logfile=/var/log/supervisor/APP_DOMAIN/QUEUE_RECONCILE.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
autostart=true
autorestart=true
user=APP_SERVER_USER
numprocs=QUEUE_NUM_PROCS

[program:QUEUE_IMAGE_EXPORT_LISTENER]
process_name=%(program_name)s_%(process_num)02d
directory=APP_CURRENT_PATH
command=php artisan queue:work beanstalkd --force --queue=QUEUE_IMAGE_EXPORT_LISTENER --tries=3
redirect_stderr=true
stdout_logfile=/var/log/supervisor/APP_DOMAIN/QUEUE_IMAGE_EXPORT_LISTENER.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
autostart=true
autorestart=true
user=APP_SERVER_USER
numprocs=QUEUE_NUM_PROCS

[program:QUEUE_RECONCILIATION_LISTENER]
process_name=%(program_name)s_%(process_num)02d
directory=APP_CURRENT_PATH
command=php artisan queue:work beanstalkd --force --queue=QUEUE_RECONCILIATION_LISTENER --tries=3
redirect_stderr=true
stdout_logfile=/var/log/supervisor/APP_DOMAIN/QUEUE_RECONCILIATION_LISTENER.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
autostart=true
autorestart=true
user=APP_SERVER_USER
numprocs=QUEUE_NUM_PROCS

[program:QUEUE_TESSERACT_OCR_LISTENER]
process_name=%(program_name)s_%(process_num)02d
directory=APP_CURRENT_PATH
command=php artisan queue:work beanstalkd --force --queue=QUEUE_TESSERACT_OCR_LISTENER --tries=3
redirect_stderr=true
stdout_logfile=/var/log/supervisor/APP_DOMAIN/QUEUE_TESSERACT_OCR_LISTENER.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
autostart=true
autorestart=true
user=APP_SERVER_USER
numprocs=QUEUE_NUM_PROCS

[program:QUEUE_WORKFLOW]
process_name=%(program_name)s_%(process_num)02d
directory=APP_CURRENT_PATH
command=php artisan queue:work beanstalkd --force --queue=QUEUE_WORKFLOW --tries=3
redirect_stderr=true
stdout_logfile=/var/log/supervisor/APP_DOMAIN/QUEUE_WORKFLOW.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
autostart=true
autorestart=true
user=APP_SERVER_USER
numprocs=QUEUE_NUM_PROCS

